# gruntworker.crontab
*/10 * * * * root docker run gruntworker 2>&1 >> /var/log/gruntworker.log
