# gruntworker.crontab
*/10 * * * * root docker run --rm=true gruntworker 2>&1 >> /var/log/gruntworker.log
